# supervisor/05-netbird.conf
[program:netbird]
command=/usr/local/bin/netbird up
autorestart=True
startsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
